# deploy/odoo-webhook.service
[Unit]
Description=Odoo GitHub Webhook Listener
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/odoo-vps
EnvironmentFile=/opt/odoo-vps/.env
ExecStart=/usr/bin/python3 /opt/odoo-vps/deploy/webhook.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
